# frontend/supervisor.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:webpack]
command=bash -c 'yarn start'